# base/files/opt/t2-startup.service
[Unit]
Description=T2 hardware setup
After=systemd-user-sessions.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/opt/scripts/t2/t2-startup.sh

[Install]
WantedBy=multi-user.target
